# roles/setup_os/usb_autorunner/files/usb-discovery.service.j2
[Unit]
Description=USB Discovery

[Service]
Type=simple
ExecStart=/bin/bash {{ usb_autorunner_install_path }}/usb-discovery.sh

[Install]
WantedBy=multi-user.target
